# fw.cron
*/1 * * * * root /bin/sh /home/fw/main.sh